# deploy/systemd/lions-creek-rewards-expire.service
# /etc/systemd/system/lions-creek-rewards-expire.service
[Unit]
Description=Lions Creek Rewards - Expire job (calls /jobs/expire)

[Service]
Type=oneshot
EnvironmentFile=/etc/lions-creek-rewards/lions-creek-rewards.env
ExecStart=/usr/bin/curl -fsS -H "x-job-token: ${JOB_TOKEN}" "http://127.0.0.1:3000/jobs/expire?shop=${PRIMARY_SHOP}"
